# config/supervisord.conf
[supervisord]
nodaemon=true

[program:xander]
command=ruby bot.rb
autorestart=true
