# config/supervisord.conf
[supervisord]
nodaemon=true

[program:xander]
command=ruby bot.rb
autorestart=true
stdout_logfile=/var/log/xander.log
